# logrotate-ntripcaster.timer
[Unit]
Description=Hourly Log Rotation for ntripcaster

[Timer]
OnCalendar=Hourly
Persistent=true

[Install]
WantedBy=timers.target

