# logrotate-ntripcaster.service
[Unit]
Description=Rotate and Compress ntripcaster Logs

[Service]
ExecStart=/usr/sbin/logrotate /usr/local/ntripcaster/ntripcaster.logrotate

